# n1/files/etc/crontabs/root
1 2 * * * $HOME/cron/daily.sh >/dev/null 2>&1
0 6 1 * * $HOME/cron/monthly.sh >/dev/null 2>&1
